# systemd/ip400_server.service
[Unit]
Description=IP400 Flask + UDP Server
After=network.target ip400_spi.service

[Service]
EnvironmentFile=-/etc/ip400/ip400.env
Environment=IP400_WEB_HOST=${IP400_SERVER_HOST}
Environment=IP400_WEB_PORT=${IP400_SERVER_PORT}
User=ip400
WorkingDirectory=${IP400_HOME}/ip400_server
ExecStart=${IP400_SERVER_VENV}/bin/python3 ${IP400_HOME}/ip400_server/ip400_server.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/ip400_server.log
StandardError=append:/var/log/ip400_server.log

[Install]
WantedBy=multi-user.target
